# crontab
* * * * * php /Users/tuan/Desktop/shopmobile/artisan schedule:run
